# mern-pipeline.service
[Unit]
Description=MERN Pipeline
After=network.target

[Service]
Type=forking
RemainAfterExit=yes
TimeoutStartUSec=10min
ExecStart=mern-pipeline up -d
ExecStop=mern-pipeline stop

[Install]
WantedBy=multi-user.target
